# Server/BackEndMonitoring.service
[Unit]
Description=Simple daemon template
After=network.target

[Service]
ExecStart=/home/Server

[Install]
WantedBy=multi-user.target
